# sites/studio/deploy/inertia-studio.service
[Unit]
Description=Inertia Studio Web Server
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/home/agent/inertia/sites/studio
EnvironmentFile=/home/agent/inertia/.env
ExecStart=/usr/bin/node dist/server/entry.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=inertia-studio

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/agent/inertia
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
